# rails_activejob_devspaces/devspaces/docker/assets/services/redis/run
#!/bin/bash

set -e
exec 2>&1

exec chpst -u redis redis-server /etc/redis/redis.conf --daemonize no
